# deployment/systemd/profiler.service
[Unit]
Description=wlan client profiler service
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/local/bin/profiler
